# desktop/agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- LaunchAgent template (§10). `patch.sh --watch` fills in the script + log paths and
     installs it to ~/Library/LaunchAgents/. User-scope, no root. WatchPaths fire when a
     Squirrel update swaps the original bundle; watch.sh then settles + re-patches. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-rtl.watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__WATCH_SH__</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Applications/Claude.app/Contents/Info.plist</string>
    <string>/Applications/Claude.app/Contents/Resources/app.asar</string>
  </array>

  <!-- Coalesce the burst of events during a bundle swap. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- launchd starts with a bare environment; give it node/npx + standard tools. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG__</string>
  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
